# supervisord.conf
[supervisord]
nodaemon=true

[program:gunicorn]
command=gunicorn -c gunicorn_config.py django_project.wsgi:application
directory=/app
autostart=true
autorestart=true
redirect_stderr=true

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
